# core/supervisor/supervisord.conf
[supervisord]
user=root
logfile=%(ENV_LOG_FOLDER_CONTAINER)s/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB                      ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10                         ; # of main logfile backups; 0 means none, default 10
loglevel=%(ENV_LOG_LEVEL)s                 ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid               ; supervisord pidfile; default supervisord.pid
nodaemon=true                              ; start in foreground if true; default false
silent=false                               ; no logs to stdout if true; default false
minfds=1024                                ; min. avail startup file descriptors; default 1024
minprocs=200                               ; min. avail process descriptors;default 200

[inet_http_server]
port=%(ENV_SERVER_URI)s:%(ENV_SUPERVISOR_PORT)s

; [program:build_worker]
; command=python -m backend.server worker
; process_name=%(program_name)s-%(process_num)s
; numprocs=%(ENV_WORKER_NUMBER)s
; directory=/home/%(ENV_DOCKER_USER)s
; stopsignal=TERM
; autostart=true
; autorestart=true

[program:build_server]
command=python -m core
process_name=%(program_name)s-%(process_num)s
numprocs=1
directory=/home/%(ENV_DOCKER_USER)s/app
stopsignal=TERM
autostart=true
autorestart=true
stdout_logfile=%(LOG_FOLDER_CONTAINER)s/server.log
